# start-gst.conf
[program:startgst]
command=/usr/local/bin/start-gst.sh
autorestart=true
user=pi
redirect_stderr=true
stdout_logfile=/var/log/gst.log
